# root/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash
set -e

umask 077
exec /usr/local/sbin/php-fpm --nodaemonize
