# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
set -e

exec nginx -g "daemon off;"
